# s6/apt-cacher-ng/run
#!/bin/sh

exec /usr/sbin/apt-cacher-ng -c /etc/apt-cacher-ng SocketPath=/var/run/apt-cacher-ng/socket pidfile=/var/run/apt-cacher-ng/pid foreground=1
